# ansible/roles/common/files/automatic_reboot.timer
[Unit]
Description=Automatic reboot

[Timer]
OnBootSec=30d
RandomizedDelaySec=24h
Unit=reboot.target

[Install]
WantedBy=timers.target